# argsms-bot.service
[Unit]
Description=ARGSMS Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/argsms
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /root/argsms/bot.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argsms-bot

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
